# crontab
# Update cargo index
30 * * * * zsh $HOME/cronjobs/update-cargo.sh
